# deploy/hireops-backend.service
[Unit]
Description=HireOps AI backend (FastAPI + uvicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=hireops
Group=hireops
WorkingDirectory=/opt/hireops/backend
EnvironmentFile=/opt/hireops/backend/.env
# Original CV binaries (direct uploads) live here, namespaced per
# tenant (tenant_<id>/candidate_<id>/...). Path is inside the
# ReadWritePaths set above so the sandboxed unit can write to it.
Environment=HIREOPS_UPLOAD_DIR=/var/lib/hireops/uploads
# Single worker on purpose: the mailbox auto-pickup listener is an in-process
# asyncio task. With multiple workers each one spawns its own listener and they
# race on the same IMAP inbox, doubling LLM classifier costs.
ExecStart=/opt/hireops/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8017 --workers 1
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/hireops/backend /var/lib/hireops
ProtectHome=true

[Install]
WantedBy=multi-user.target
